# supervisor/super.conf
[unix_http_server]
file = /home/vunm/var/log/supervisor_telegram.sock

[supervisord]
directory = /home/vunm
nodaemon  = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = /home/vunm/var/log/supervisor_telegram.sock

[program:api_8101]
command = python /home/vunm/telegramBot/StudyEnglish/detailed_example.py
autorestart = true